# nginx-php/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:php5-fpm]
command=/usr/sbin/php5-fpm --allow-to-run-as-root -y /etc/php5/fpm/php-fpm.conf -c /etc/php5/fpm/php.ini -F
